# tz-php7/config/supervisor.conf
[supervisord]
nodaemon=true
logfile=/runtime/log/supervisor.log

[program:nginx]
command=nginx -g "daemon off;"

[program:php-fpm]
command=/usr/sbin/php-fpm7 --nodaemonize
