# tz-py2-uwsgi/config/supervisor.conf
[supervisord]
nodaemon=true
logfile=/runtime/log/supervisor.log

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"

[program:uwsgi]
command=/usr/sbin/uwsgi /runtime/config/uwsgi.ini
